# install/etc/init/robot.conf
env ROS_LOG_DIR=/var/log/ros

start on roscore_is_up
stop on roscore_is_down

respawn

script
    exec su ros -c ". /etc/ros/melodic/setup.sh && roslaunch /etc/ros/melodic/robot.launch"
end script
