# install/etc/init/sixad.conf
start on (started dbus and bluetooth-device-added)
stop on stopping dbus

respawn

exec /usr/bin/sixad --start
