# launch/install/etc/init/ps3joy.conf
start on runlevel [2345]

respawn

exec /opt/ros/hydro/lib/ps3joy/ps3joy.py --continuous-output
